# launchd/Library/LaunchAgents/com.jenner.tmux-claude-titles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jenner.tmux-claude-titles</string>
    <!-- Runs ~/bin/tmux-claude-titles-loop.sh, which loops forever and no-ops while no
         tmux server is up. launchd owns its lifecycle instead of tmux's `run-shell -b`. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jenner/bin/tmux-claude-titles-loop.sh</string>
    </array>
    <!-- launchd hands agents a bare PATH; the loop and its workers shell out to `tmux`
         (Homebrew) plus the usual coreutils. LANG forces a UTF-8 locale: without it the
         agent runs in the C locale, where tmux rewrites the TAB separators in our `-F`
         format strings as '_' and the worker's IFS=tab parsing skips every pane. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/jenner/Library/Logs/tmux-claude-titles.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jenner/Library/Logs/tmux-claude-titles.log</string>
</dict>
</plist>
